# systemd/her-axera-axllm.service.in
[Unit]
Description=Her.axera AX LLM inference service (ax-llm)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=axllm serve --port @AXLLM_PORT@ --model-path @MODEL_PATH@
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
